# scripts/deskos-news-pipeline.timer
[Unit]
Description=Run DeskOS News Intelligence Pipeline every 3 days

[Timer]
OnBootSec=10min
OnUnitActiveSec=3d
Persistent=true

[Install]
WantedBy=timers.target
