# scripts/deskos-kiosk.service
[Unit]
Description=DeskOS Kiosk (Electron)
After=graphical.target deskos-backend.service
Wants=deskos-backend.service
Requires=graphical.target

[Service]
Type=simple
User=__DESKOS_USER__
Environment=DISPLAY=:0
Environment=XAUTHORITY=__DESKOS_HOME__/.Xauthority
ExecStartPre=/bin/sleep 3
ExecStartPre=__DESKOS_REPO_DIR__/scripts/apply-orientation.sh
ExecStart=__DESKOS_REPO_DIR__/scripts/launch-kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
