# scripts/deskos-backend.service
[Unit]
Description=DeskOS Backend
After=network.target
Wants=network.target

[Service]
Type=simple
User=__DESKOS_USER__
WorkingDirectory=__DESKOS_REPO_DIR__
ExecStart=/usr/bin/npm run start --workspace=packages/backend
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
