# scripts/deskos-kiosk.service
[Unit]
Description=DeskOS Kiosk (Chromium)
After=graphical.target deskos-backend.service
Wants=deskos-backend.service
Requires=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 3
ExecStart=/home/pi/rpi-home-assistant/scripts/launch-kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
